# sys-units/anonsurfd.service
[Unit]
Description="Anonsurf Daemon Service"
Before=poweroff.target reboot.target
After=network.target

[Service]
# ExecStartPre=/etc/anonsurf/make-torrc
ExecStart=/etc/anonsurf/anondaemon start
ExecReload=/etc/anonsurf/anondaemon restart
ExecStop=/etc/anonsurf/anondaemon stop
RemainAfterExit=yes

# Hardening
AppArmorProfile=-system_tor
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=full
ReadOnlyDirectories=/

[Install]
WantedBy=default.target